# search_jobs/etc/crontab
* * * * * echo hello  >> /logs/job.output 2>> /logs/job.errors

